# deployment/fintube-scout.service
[Unit]
Description=FinTube scout — discover & distill new relevant YouTube videos, push to Telegram
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=user
# The scout runs IN the backend container (where yt-dlp, Redis and the Opus worker bus
# live). We just trigger the in-app background task over the loopback API; the backend
# does discovery → distill → Telegram itself. AGENT_AUTH_DISABLED=true on the box, so no
# auth header is needed; nginx/Cloudflare never see this loopback call.
ExecStart=/usr/bin/curl --fail --silent --show-error --max-time 30 \
    --retry 3 --retry-delay 10 --retry-connrefused \
    -X POST http://127.0.0.1:8000/api/fintube/scout \
    -H "Content-Type: application/json" -d "{}"
TimeoutStartSec=120
