# worker/sector-rotation.timer
# systemd timer: fire the sector-rotation warm once per day, after crack-a-dawn.
#
# 6:15 AM PT — a quarter hour behind the crack-a-dawn briefing so the box isn't
# doing both heavy sweeps at once. The schedule is expressed in America/Los_Angeles
# so it tracks PT across DST without edits.
#
# Persistent=true => if the box was asleep/off at 6:15, the warm runs on the next
# boot rather than silently skipping the day.

[Unit]
Description=Daily warm of the sector-rotation snapshot
Requires=sector-rotation.service

[Timer]
OnCalendar=*-*-* 06:15:00 America/Los_Angeles
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
